# deploy/web/systemd/agent-news-web-api.service
[Unit]
Description=Agent News dashboard API (FastAPI / uvicorn)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=user1
Group=user1
WorkingDirectory=/home/user1/gosb_bot
Environment=PYTHONPATH=/home/user1/gosb_bot/src
Environment=AGENT_NEWS_WEB_HOST=127.0.0.1
Environment=AGENT_NEWS_WEB_PORT=8765
# Uncomment to enforce Basic auth at the API layer (in addition to nginx):
# Environment=AGENT_NEWS_WEB_USER=radar
# Environment=AGENT_NEWS_WEB_PASSWORD=change-me

# OpenClaw agent integration (for /api/v1/chat endpoint)
Environment=OPENCLAW_BIN=/home/user1/.npm-global/bin/openclaw
Environment=OPENCLAW_MODEL=openai/gpt-5.5
Environment=OPENCLAW_TIMEOUT=240

# Use the project's venv if you have one. Plain python3 also works.
ExecStart=/usr/bin/python3 -m agent_news.web

Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
